# infra/packer/service_files/bento-executor.service
[Unit]
Description=Boundless Bento Prover Service
After=network.target

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/etc/environment
WorkingDirectory=/opt/boundless
ExecStart=/usr/local/bin/agent -t exec
Restart=always
RestartSec=10
TimeoutStartSec=1800
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
